# conf/database
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

#extra_commands=""
#extra_started_commands=""

depend() {
	use logger
	need net
}

PIDFILE="/var/run/${SVCNAME}.pid"
DIR="/home/${SVCNAME}"
BINARY="$DIR/Database"
ARGS="-d"
USER="${SVCNAME}:${SVCNAME}"

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start  \
	    --interpreted --make-pidfile --pidfile "${PIDFILE}" --user "${USER}" \
            --chdir "${DIR}" --background --stdout /home/database/log/service.log \
            --exec "${BINARY}" -- ${ARGS}
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --exec "${BINARY}" \
	    --interpreted --pidfile "${PIDFILE}" --user "${USER}" --quiet
	eend $?
}
