# services/mr/mr
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

#extra_commands=""
#extra_started_commands=""

depend() {
	use logger
	need net
}

SVCNAME=mr
PIDFILE="/var/run/${SVCNAME}.pid"
DIR="/home/${SVCNAME}"
START="-sname mr -setcookie mr_cookie -noshell -pa /home/mr/ebin /home/mr/deps/*/ebin /home/mr/priv/code -s inets -s yaws -s ibrowse"
BINARY="/usr/bin/erl"
USER="${SVCNAME}:${SVCNAME}"

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start \
	    --interpreted --make-pidfile --pidfile "${PIDFILE}" --user "${USER}" \
            --chdir "${DIR}" --background --stdout /home/mr/service.log --exec "${BINARY}" -- ${START} 
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop \
	    --interpreted --pidfile "${PIDFILE}" --user "${USER}" --quiet  --exec "${BINARY}"
	eend $?
}
